# board/batocera/fsoverlay/etc/profile.d/30-welcome.sh
# Add KNULLI logo and some alias, sourcing of $HOME/.bashrc can be added to $HOME/.profile
echo 
echo " _  ___   _ _   _ _     _     ___ "
echo "| |/ / \ | | | | | |   | |   |_ _|"
echo "| ' /|  \| | | | | |   | |    | | "
echo "| . \| |\  | |_| | |___| |___ | | "
echo "|_|\_\_| \_|\___/|_____|_____|___|"
echo
batocera-info 2>/dev/null
echo "OS version: $(batocera-version)"
echo
