# board/batocera/rockchip/rk3566/fsoverlay/etc/udev/rules.d/60-hdmi.rules
ACTION=="change",SUBSYSTEM=="drm",ENV{DEVPATH}=="/devices/platform/display-subsystem/drm/card0",RUN+="/usr/bin/hdmi_audio_switch.sh"


